# tools/youtube-upload/com.giis.youtube-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  GIIS YouTube daily uploader · launchd job

  WHAT IT DOES
    Runs tools/youtube-upload/daily.sh every day at 09:00 local Mac time.
    daily.sh prints the queue status, then uploads up to 4 pending lessons
    (stays under YouTube Data API v3 quota).

  WHY 09:00
    YouTube Data API quota resets at midnight Pacific Time. 09:00 local time
    on the user's Mac is comfortably after that for any US time zone, while
    still letting the run finish before the workday's first cup of coffee.
    If you're on a non-Pacific Mac and want exact PT timing, change Hour to:
      ET / EDT     → 12 (noon ET)
      CT / CDT     → 11
      MT / MDT     → 10
      AKT / AKDT   → 08
      HAT          → 06
    Or to follow PT regardless of local zone, switch this from
    StartCalendarInterval to a Program that calls `sleep $(seconds_until_9am_PT)`.

  INSTALL
    cp tools/youtube-upload/com.giis.youtube-daily.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.giis.youtube-daily.plist

  CHECK
    launchctl list | grep giis              # should show the job
    tail -f ~/Library/Logs/giis-youtube-daily.log

  STOP
    launchctl unload ~/Library/LaunchAgents/com.giis.youtube-daily.plist
    rm ~/Library/LaunchAgents/com.giis.youtube-daily.plist

  RUN ONCE NOW (without waiting for the 9am tick)
    launchctl kickstart -k gui/$(id -u)/com.giis.youtube-daily

  IMPORTANT
    Edit the WorkingDirectory below if your repo lives anywhere other than
    /Users/alanhdchu/giis-website.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.giis.youtube-daily</string>

    <!-- Paused by Alan on 2026-05-22: do not auto-upload GIIS lesson videos. -->
    <key>Disabled</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/alanhdchu/giis-website</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alanhdchu/giis-website/tools/youtube-upload/daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the Mac is asleep at 09:00, run as soon as it wakes. -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Don't hold the system busy — launchd should batch us with other low-priority jobs. -->
    <key>LowPriorityIO</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/alanhdchu/Library/Logs/giis-youtube-daily.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alanhdchu/Library/Logs/giis-youtube-daily.err.log</string>
</dict>
</plist>
